# sleepybox.service
[Unit]
Description=SleepyBox

[Service]
ExecStartPre=killall sleepybox
ExecStart=/usr/bin/sleepybox

[Install]
WantedBy=multi-user.target
